# templates/etc/systemd/system/containerd.service.d/override.conf.j2
{{ ansible_managed | comment }}
[Service]
ExecStartPre=
